# Prototype/session_messenger_bridge/session_bot_secretarius.service
[Unit]
Description=Secretarius Session Messenger Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/mauceric/Secretarius/Prototype/session_messenger_bridge

Environment=HOME=/home/mauceric
Environment=PATH=/home/mauceric/.nvm/versions/node/v24.13.0/bin:/home/mauceric/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=SESSION_MESSENGER_API_URL=http://127.0.0.1:8002
Environment=SESSION_MESSENGER_JOURNAL_PATH=/home/mauceric/Secretarius/Prototype/logs/session_messenger.log
Environment=SESSION_MESSENGER_SQLITE_PATH=/home/mauceric/Secretarius/Prototype/session_messenger_bridge/session_messenger_messages.sqlite
Environment=SESSION_BOT_DISPLAY_NAME=Secretarius

ExecStart=/home/mauceric/.bun/bin/bun run start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
