# Prototype/session_messenger_bridge/session-messenger-bot.service
[Service]
WorkingDirectory=/home/mauceric/Secretarius/Prototype/session_messenger_bridge

Environment=HOME=/home/mauceric
Environment=PATH=/home/mauceric/.nvm/versions/node/v24.13.0/bin:/home/mauceric/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=SESSION_MESSENGER_API_URL=http://127.0.0.1:8002
Environment=SESSION_MESSENGER_JOURNAL_PATH=/home/mauceric/Secretarius/Prototype/logs/session_messenger.log
Environment=SESSION_MESSENGER_SQLITE_PATH=/home/mauceric/Secretarius/Prototype/session_messenger_bridge/session_messenger_messages.sqlite

ExecStart=/home/mauceric/.bun/bin/bun run start
