# deploy/systemd-user/open-webui.service
[Unit]
Description=Open WebUI
After=network-online.target secretarius_ollama.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/mauceric/Secretarius
Environment=PATH=/home/mauceric/Secretarius/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/mauceric/Secretarius/.venv/bin/open-webui serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
