# deploy/systemd-user/secretarius-memory.service
[Unit]
Description=Secretarius Memory API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/mauceric/Secretarius
EnvironmentFile=%h/.config/secretarius/memory.env
ExecStart=/home/mauceric/Secretarius/.venv/bin/python /home/mauceric/Secretarius/serverMemory.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target

