# deploy/systemd-user/secretarius_ollama.service
[Unit]
Description=Secretarius Ollama service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/mauceric/Secretarius
Environment=PATH=/home/mauceric/Secretarius/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=HSA_OVERRIDE_GFX_VERSION=10.3.0
Environment=OLLAMA_HOST=0.0.0.0:11434
Environment=OLLAMA_CONTEXT_LENGTH=8192
ExecStart=/usr/local/bin/ollama serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
