# deploy/systemd-user/secretarius-telegram.service
[Unit]
Description=Secretarius Telegram Adapter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/mauceric/Secretarius
EnvironmentFile=%h/.config/secretarius/telegram.env
ExecStart=/home/mauceric/Secretarius/.venv/bin/python /home/mauceric/Secretarius/serverTelegram.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target

