# deploy/systemd-user/secretarius-email.service
[Unit]
Description=Secretarius Email Adapter (IMAP/SMTP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/mauceric/Secretarius
EnvironmentFile=%h/.config/secretarius/email.env
ExecStart=/home/mauceric/Secretarius/.venv/bin/python /home/mauceric/Secretarius/serverEmail.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target

