# deploy/systemd-user/secretarius-prototype-session-snapshot.timer
[Unit]
Description=Periodic Secretarius Prototype session snapshot

[Timer]
OnBootSec=10m
OnUnitActiveSec=30m
Unit=secretarius-prototype-session-snapshot.service
Persistent=true

[Install]
WantedBy=timers.target
