# deploy/systemd-user/secretarius-prototype-session-snapshot.service
[Unit]
Description=Secretarius Prototype automatic session snapshot

[Service]
Type=oneshot
WorkingDirectory=/home/mauceric/Secretarius/Prototype
ExecStart=/home/mauceric/Secretarius/Prototype/scripts/snapshot_auto.sh
