# deploy/systemd-user/secretarius-openwebui-api.service
[Unit]
Description=Secretarius OpenAI API for Open WebUI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/mauceric/Secretarius
EnvironmentFile=%h/.config/secretarius/openwebui-api.env
ExecStart=/home/mauceric/Secretarius/.venv/bin/python /home/mauceric/Secretarius/serverOpenAI.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target

