# openclaw-config/gog-mcp.service
[Unit]
Description=gog MCP Server (SSE port 8902)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Secretarius
EnvironmentFile=%h/.openclaw/gateway.systemd.env
ExecStart=%h/Secretarius/Wiki_LM/.venv/bin/python3 %h/Secretarius/gog_mcp/mcp_server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gog-mcp

[Install]
WantedBy=default.target
