# openclaw-config/openclaw-gateway-slm.service
[Unit]
Description=OpenClaw Gateway SLM (dev, instance isolée 5.12)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.openclaw-slm/gateway.systemd.env
ExecStart=${OPENCLAW_BIN} gateway run --profile slm
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
