# openclaw-config/openclaw-injection-guard.service
[Unit]
Description=Scout Injection Guard — détection d'injection de prompt
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=%h/Secretarius/Wiki_LM/.venv/bin/python3 %h/.local/bin/injection_guard.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=injection-guard

[Install]
WantedBy=default.target
