# openclaw-config/wiki-lm-mcp.service
[Unit]
Description=Wiki_LM MCP Server (SSE port 8901)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Secretarius/Wiki_LM
ExecStart=%h/Secretarius/Wiki_LM/.venv/bin/python3 %h/Secretarius/Wiki_LM/tools/mcp_server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wiki-lm-mcp

[Install]
WantedBy=default.target
