# openclaw-config/wiki-lm-embed.timer
[Unit]
Description=Rafraîchissement périodique des embeddings Wiki_LM (toutes les 6 h)

[Timer]
OnCalendar=*-*-* 00,06,12,18:00:00
Persistent=true

[Install]
WantedBy=timers.target
